# dockerfiles/elk/files/supervisor/httpd.conf
[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND
autostart=true
log_stdout=true
log_stderr=true
