# dockerfiles/elk/files/supervisor/logstash.conf
[program:logstash]
command=/opt/logstash/bin/logstash agent -f /etc/logstash/logstash.conf
autostart=true
log_stdout=true
log_stderr=true


