# dockerfiles/shinken-enterprise-2.0.3-update/supervisor/graphite.conf
[program:graphite-web]
user = www-data
directory = /opt/graphite/webapp
environment = PYTHONPATH='/var/lib/graphite/webapp'
command = /usr/bin/gunicorn_django -b127.0.0.1:8000 -w2 graphite/settings.py
redirect_stderr = true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
startsecs=60

[program:carbon-cache]
user = www-data
command = /opt/graphite/bin/carbon-cache.py --debug start
redirect_stderr = true
startsecs=60