# dockerfiles/shinken-enterprise-2.0.3-update/supervisor/mongodb.conf
[program:shinken-mongodb]
command=/usr/bin/mongod -f /etc/mongod.conf
autostart=true
user=mongod
redirect_stderr=true
stdout_logfile=/var/log/mongod.log


