# dockerfiles/centreon/apache.conf
[program:apache]
command=/usr/sbin/httpd -DFOREGROUND
autostart=true
redirect_stderr=true