# dockerfiles/centreon/centengine.conf
[program:centengine]
command=/usr/sbin/centengine /etc/centreon-engine/centengine.cfg
priority=70
user=centreon-engine
autorestart=true
redirect_stderr=true