# dockerfiles/centreon/crond.conf
[program:crond]
command = /usr/sbin/crond -n
autorestart = true
redirect_stderr=true