# dockerfiles/centreon/mysql.conf
[program:mysql]
command=/bin/sh /usr/bin/mysqld_safe --datadir=/var/lib/mysql
autostart=true
redirect_stderr=true

