# dockerfiles/centreon/rsyslogd.conf
[program:rsyslogd]
command=rsyslogd -n
autorestart=true
redirect_stderr=true