# dockerfiles/nagvis/supervisor/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D -e
priority=10
redirect_stderr=true
startsecs=5
autorestart=true
redirect_stderr=true