# dockerfiles/thruk-shinken-2.2/apache.conf
[program:apache]
command=/usr/local/bin/startapache.sh
autostart=true
redirect_stderr=true