# dockerfiles/graphite/supervisord/carbon-cache.conf
[program:carbon-cache]
user = www-data
command = /opt/graphite/bin/carbon-cache.py --debug start
redirect_stderr = true
stdout_logfile = /var/log/supervisor/%(program_name)s.log