# dockerfiles/shinken-allinone-1.4/files/supervisor/snmpd.conf
[program:snmpd]
command = snmpd -f -Lsd -Lf /var/log/snmpd.log -p /var/run/snmpd -a
autorestart = true
log_stdout=true
log_stderr=true

