# dockerfiles/shinken-1.4/files/supervisor/sshd.conf
[program:sshd]
command = /usr/sbin/sshd -D 
autorestart = true
log_stdout=true
log_stderr=true
