# dockerfiles/shinken-arbiter-1.4/files/supervisor/shinken.conf
[program:arbiter]
command=/usr/bin/shinken-arbiter -c /etc/shinken/nagios.cfg -c /etc/shinken/shinken-specific.cfg
autostart=true
user=shinken
log_stdout=true
log_stderr=true