# dockerfiles/thruk-shinken-2.2/supervisor/rsyslog.conf
[program:rsyslog]
command=/usr/sbin/rsyslogd -n
autostart=true
redirect_stderr=true